# etc/sv/git-sunc/run
#!/usr/bin/env zsh
exec snooze -M/15 git-sunc
